# build_ramdisk/obstool-image/template-obshx1-runled/DEBIAN/preinst
#!/bin/sh
# preinst script for openblocks family
#
# see: dh_installdeb(1)

#set -e

. /etc/default/openblocks

case "$1" in
install|upgrade)
	[ -f /.chroot ] && exit 0
	if systemctl is-active runled; then
		systemctl stop runled && systemctl disable runled
	fi
	;;
abort-upgrade)
	;;
*)
	echo "preinst called with unknown argument \`$1'" >&2
	exit 1
	;;
esac



# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.

#DEBHELPER#

exit 0
