# build_ramdisk/obstool-image/template-flashcfg/lib/systemd/system/instfirm.service
[Unit]
Description=Restore factry image Openblocks family

[Service]
ExecStart=/usr/sbin/instfirm.sh
Type=oneshot

[Install]
WantedBy=multi-user.target
