# build_ramdisk/etc.squeeze/exclude.list
etc/default/openblocks
etc/exclude.list
etc/init.d/openblocks-setup
etc/init.d/pshd
etc/init.d/runled
etc/init.d/umountfs
etc/mtab
etc/udev/rules.d/70-persistent-net.rules
usr/lib/jre
usr/sbin/flashcfg
usr/sbin/flashcfg-debian
var/cache/apt
var/lib/apt/lists
var/lib/dpkg/updates
var/log
var/run
