# build_ramdisk/obstool-image/template-pshd/DEBIAN/postinst
#!/bin/sh
# preinst script for openblocks family
#
# see: dh_installdeb(1)

#set -e

. /etc/default/openblocks

case "$1" in
configure)
	[ -f /.chroot ] && exit 0
	if systemctl is-active pshd; then
		systemctl restart pshd
	else
		systemctl enable pshd && systemctl start pshd
	fi
	;;
abort-upgrade|abort-remove|abort-deconfigure)
	;;
*)
	echo "postinst called with unknown argument \`$1'" >&2
	exit 1
	;;
esac



# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.

#DEBHELPER#

exit 0
