# build_ramdisk/obstool-image/template-obs-nicled/lib/systemd/system/obs-nicled.service
[Unit]
Description=Setup NIC LED for Openblocks family
After=networking.service

[Service]
ExecStart=/usr/sbin/obs-nicled.sh
Type=oneshot

[Install]
WantedBy=multi-user.target
