# build_ramdisk/obstool-image/template-nitz/lib/systemd/system/nitz.service
[Unit]
Description=Network Identity and Time Zone for Openblocks family
After=ssh.service

[Service]
ExecStart=/usr/sbin/nitz.sh
Type=oneshot

[Install]
WantedBy=multi-user.target
