# build_ramdisk/obstool-image/template-disable-modem/DEBIAN/postinst
#!/bin/sh
# preinst script for openblocks family
#
# see: dh_installdeb(1)

#set -e

. /etc/default/openblocks

case "$1" in
configure)
	[ -f /.chroot ] && exit 0
	systemctl enable disable-modem
	;;
abort-upgrade|abort-remove|abort-deconfigure)
	;;
*)
	echo "postinst called with unknown argument \`$1'" >&2
	exit 1
	;;
esac



# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.

#DEBHELPER#

exit 0
