# build_ramdisk/uboot-image/template/DEBIAN/postinst
#!/bin/sh
# preinst script for easyblockssync
#
# see: dh_installdeb(1)

#set -e

. /etc/default/openblocks

case "$1" in
configure)
	[ -f /.chroot ] && exit 0
	if [ -e /etc/init.d/easyblockssync ]; then
		/usr/sbin/update-rc.d u-boot defaults
		/etc/init.d/u-boot start
	elif [ -f /etc/mtd0.img.xz ]; then
		if /etc/init.d/u-boot start; then
			rm -f /etc/mtd0.img.xz
			rm -f /etc/mtd1.img.xz
			rm -f /etc/mtd1_java.img.xz
		else
			exit 1
		fi
	fi
;;

abort-upgrade|abort-remove|abort-deconfigure)
;;

*)
echo "postinst called with unknown argument \`$1'" >&2
exit 1
;;
esac



# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.

#DEBHELPER#

exit 0
