# build_ramdisk/obstool-image/template-disable-modem/lib/systemd/system/disable-modem.service
[Unit]
Description=Poweroff Modem when restarting for OpenBlocks family

[Service]
ExecStart=/usr/sbin/disable-modem.sh
Type=oneshot

[Install]
WantedBy=reboot.target
